# aurora-search.service
[Unit]
Description=Aurora Search Backend (ddgs API)
After=network.target

[Service]
ExecStart=%h/.local/share/pipx/venvs/ddgs/bin/ddgs api --host 127.0.0.1 --port 8888
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
